# main/files/renewal/renewssl.service
[Unit]
Description=Renews letencrypt SSL certificates

[Service]
Type=simple
ExecStart=/bin/sh -c 'letsencrypt renew'
